# examples/runit/log/run
#!/bin/sh
exec svlogd -tt /var/log/sidekiq
